# resources/fty-metric-ambient-location.service.in
[Unit]
Description=fty-metric-ambient-location service
After=network.target malamute.service
Requires=malamute.service
PartOf=bios.target

After=fty-asset.service
Requires=fty-asset.service
BindsTo=fty-asset.service

[Service]
Type=simple
User=@MAL_USER@
# User=@uid@
EnvironmentFile=-/usr/share/fty/etc/default/fty
EnvironmentFile=-/usr/share/fty/etc/default/fty__%n.conf
Environment="prefix=/usr"
Environment='SYSTEMD_UNIT_FULLNAME=%n'
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/@PROJECT_NAME@
Restart=always

[Install]
WantedBy=bios.target
